# deploy/systemd/izotop-connect-bot.service
[Unit]
Description=Izotop Connect Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/izotop-connect-bot
EnvironmentFile=/opt/izotop-connect-bot/.env
ExecStart=/opt/izotop-connect-bot/.venv/bin/python -m izotop_connect_bot
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

